# broker/etc/init.d/terrapi
#!/sbin/openrc-run

command="/usr/bin/terrapi"
command_args="/etc/terrapi/env"
command_background=true
command_user="terra"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
        need localmount net
}

