# broker/etc/init.d/terrapi-broker
#!/sbin/openrc-run

command="/usr/bin/terrapi-broker"
command_args="--env-file /etc/terrapi-broker/config"
command_background=true
command_user="terra"
output_log="/var/log/terrapi-broker.log"
error_log="/var/log/terrapi-broker.log"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
        need localmount net
}

