# broker/etc/init.d/terrapi-broker
#!/sbin/openrc-run

command="/usr/bin/terrapi-broker"
command_args="/etc/terrapi-broker/config"
command_background=true
command_user="terra"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
        need localmount net
}

